# deploy/systemd/smarttour-nginx-host-report.service
[Unit]
Description=SmartTour Nginx host and bot traffic report
After=docker.service

[Service]
Type=oneshot
UMask=0027
EnvironmentFile=-/etc/default/smarttour-ops
WorkingDirectory=/opt/smarttour
ExecStart=/opt/smarttour/scripts/nginx-host-report.sh
TimeoutStartSec=2min
StandardOutput=append:/var/log/smarttour/nginx-host-report.log
StandardError=append:/var/log/smarttour/nginx-host-report.log
