# deploy/systemd/smarttour-restore-drill.service
[Unit]
Description=SmartTour weekly PostgreSQL restore drill
After=docker.service smarttour-postgres-backup.service

[Service]
Type=oneshot
UMask=0027
EnvironmentFile=-/etc/default/smarttour-ops
WorkingDirectory=/opt/smarttour
ExecStart=/opt/smarttour/scripts/restore-drill-postgres.sh
TimeoutStartSec=120min
StandardOutput=append:/var/log/smarttour/restore-drill.log
StandardError=append:/var/log/smarttour/restore-drill.log
