# deploy/systemd/smarttour-postgres-backup.timer
[Unit]
Description=Run SmartTour PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 02:15:00
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
